# etc/systemd/system/nvidia-cdi-generate.timer
[Unit]
Description=Regenerate NVIDIA CDI spec daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
