# etc/systemd/system/nvidia-cdi-generate.timer
[Unit]                                                           
Description=Regenerate NVIDIA CDI spec on boot and daily
                                                          
[Timer]                                                                                                                                                                
OnBootSec=0
OnCalendar=daily                                                                                                                                                       
Persistent=true 
                 
[Install]
WantedBy=timers.target
